# .config/autostart/dunst.desktop
[Desktop Entry]
Name=Start Dunst 
Exec=dunst
Terminal=false
Type=Application
