# .config/autostart/waybar.desktop
[Desktop Entry]
Name=Waybar
Exec=waybar
Terminal=false
Type=Application
